# blog.upstart
description "apostopher blog"
author      "rahul devaskar"

start on startup
stop on shutdown

respawn
respawn limit 10 10

env NODE_ENV=production

script
    chdir /home/apos/website/blog
    exec npm start >> /var/log/blog-apostopher.sys.log 2>&1
end script